# Introduction To Linux (LFS101) /projects/systemd/healthcheck.timer
[Unit]
Description=Hourly health check timer


[Timer]
OnCalendar=hourly
Persistent=true


[Install]
WantedBy=timers.target
